# deploy/systemd/normalized/syndicate-discord.service
[Unit]
Description=Syndicate - Discord Content Router Bot
After=network.target

[Service]
Type=simple
User=ali_shakil_backup_gmail_com
WorkingDirectory=/home/ali_shakil_backup_gmail_com/syndicate
Environment="PYTHONPATH=/home/ali_shakil_backup_gmail_com/syndicate" # pragma: allowlist secret
ExecStart=/home/ali_shakil_backup_gmail_com/miniforge3/envs/syndicate/bin/python -m src.digest_bot.discord.bot
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
